# deploy/systemd/yts-render-hub-reload.service.in
[Unit]
Description=Restart YTS Render review hub after repository changes
After=yts-render-hub.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart yts-render-hub.service
